# pusher.dev.conf
[program:pusher]
command=php bin/pusher.php
directory=/Users/james/code/silex-test
user=nobody
autostart=false
autorestart=true
redirect_stderr=True