# pusher.conf
[program:pusher]
command=php bin/pusher.php
directory=/home/web/sites/silex-test.jmoz.co.uk
user=nobody
autostart=true
autorestart=true
redirect_stderr=True