# pusher.conf
[program:pusher]
command=php bin/pusher.php
directory=/home/web/sites/silex-test
user=nobody
autostart=false
autorestart=true
redirect_stderr=True